# dashcam.service
[Unit]
Description=Dashcam Startup
After=user.target

[Service]
